# php/supervisord.ini
[supervisord]
nodaemon=true
user=root

[program:queue]
command=php /var/www/artisan horizon
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/horizon.log
stdout_logfile_maxbytes=0
stopwaitsecs=3600

[program:schedule]
command=php /var/www/artisan schedule:run
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/schedule.log
stdout_logfile_maxbytes=0
stopwaitsecs=3600
